# everest-github-scholar-google/src/scholar/server.conf
description "Scholar Service"
author  "Raúl Ortega <rortega@expertsystem.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir /home/everest/scholar
exec python3.4 server.py
respawn
